# deploy/systemd/ec-ai-wallet-preview.service
[Unit]
Description=EC AI Wallet Preview
After=network-online.target apppack-view-service-api.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/src/ec-ai-wallet
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 4173
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
